# templates/labmd.service
[Unit]
Description=LabMD - Lab Monitoring and Documentation
Documentation=https://github.com/SheepTAO/labmd
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/labmd server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=labmd

# Security hardening (limited due to need for system monitoring)
# Note: Requires root access to read all users' home directories and system stats
PrivateTmp=true
ReadWritePaths=/etc/labmd

[Install]
WantedBy=multi-user.target
